# templates/default/etc/init.d/knotx.erb
#!/bin/sh

#
# chkconfig: - 80 20
# processname: <%= @knotx_id %>
# description: Knotx - reactive microservice assembler
#

### BEGIN INIT INFO
# Provides: <%= @knotx_id %>
# Required-Start: $network $syslog
# Required-Stop: $network $syslog
# Default-Start:
# Default-Stop:
# Short-Description: start and stop <%= @knotx_id %>
### END INIT INFO

. /etc/rc.d/init.d/functions

LOCK_FILE="/var/lock/subsys/<%= @knotx_id %>"
JAR_FILE="<%= @knotx_root_dir %>/knotx.jar"

CONFIG_FILE="<%= @knotx_root_dir %>/knotx.conf"

[ -e $CONFIG_FILE ] && . $CONFIG_FILE

# Removing endlines from variable for loading
KNOTX_OPTS=$(echo $KNOTX_OPTS|tr -d '\n')
KNOTX_CONFIG=$(echo $KNOTX_CONFIG|tr -d '\n')

do_start()
{
  if [ ! -f "$LOCK_FILE" ] ; then
    echo -n $"Starting <%= @knotx_id %>: "
    runuser -l "<%= @knotx_user %>" -c "java $KNOTX_OPTS -jar $JAR_FILE $KNOTX_CONFIG >> <%= @knotx_log_dir %>/<%= @knotx_id %>.log 2>&1 &" && echo_success || echo_failure
    RETVAL=$?
    echo
    if [ $RETVAL -eq 0 ] ; then
      touch $LOCK_FILE
    else
      echo "Error while starting <%= @knotx_id %>. Error code from java: [$RETVAL]"
    fi
  else
    echo "<%= @knotx_id %> is locked. Lock file [$LOCK_FILE] is present therefore instance is probably running."
    RETVAL=0
  fi
}
do_stop()
{
  echo -n $"Stopping <%= @knotx_id %>: "
  pid=`ps -aefw | grep "jar $JAR_FILE" | grep -v " grep " | awk '{print $2}'`
  kill -9 $pid > /dev/null 2>&1 && echo_success || echo_failure
  RETVAL=$?
  echo
  [ $RETVAL -eq 0 ] && rm -f $LOCK_FILE
}
do_status()
{
  if [ -f "$LOCK_FILE" ] ; then
    echo "<%= @knotx_id %> is running..."
    exit 0
  else
    echo "<%= @knotx_id %> is stopped"
    exit 3
  fi
}

case "$1" in
start)
  do_start
;;
stop)
  do_stop
;;
restart)
  do_stop
  do_start
;;
status)
  do_status
;;
*)
echo "Usage: $0 {start|stop|restart}"
RETVAL=1
esac

exit $RETVAL
