# templates/default/etc/systemd/system/knotx.service.erb
[Unit]
Description=<%= @id %>
After=network.target

[Service]
EnvironmentFile=<%= @home_dir %>/knotx.conf
ExecStart=<%= @java_home %>/bin/java \
  ${KNOTX_GENERAL_OPTS} \
  ${KNOTX_GC_OPTS} \
  ${KNOTX_JMX_OPTS} \
  ${KNOTX_DEBUG_OPTS} \
  -Xms${KNOTX_MIN_HEAP}m \
  -Xmx${KNOTX_MAX_HEAP}m \
  -XX:+UseBiasedLocking \
  -XX:BiasedLockingStartupDelay=0 \
  -Dlogback.configurationFile=<%= ::File.join(@conf_dir, 'logback.xml') %> \
  -Dvertx.logger-delegate-factory-class-name=io.vertx.core.logging.SLF4JLogDelegateFactory \
  -Dknotx.home=<%= @home_dir %> \
  -Dvertx.cacheDirBase=<%= ::File.join(@home_dir, '.vertx' ) %> \
  -Dvertx.cli.usage.prefix=knotx \
  -classpath '<%= @conf_dir %>:<%= ::File.join(@lib_dir, '*') %>' \
  ${KNOTX_EXTRA_OPTS} \
  io.vertx.core.Launcher run-knotx
Type=simple
PIDFile=/var/run/<%= @id %>
User=<%= @user %>
LimitNOFILE=<%= @open_file_limit %>
SuccessExitStatus=143

[Install]
WantedBy=default.target
