# templates/default/etc/systemd/system/knotx.service.erb
[Unit]
Description=<%= @knotx_id %>
After=network.target

[Service]
EnvironmentFile=<%= @knotx_root_dir %>/knotx.conf
ExecStart=/bin/sh -c '/usr/bin/java \
  -Dlogback.configurationFile=<%= @knotx_root_dir %>/logback.xml \
  -Dvertx.cacheDirBase=<%= @knotx_root_dir %>/.vertx \
  -Xms${KNOTX_MIN_HEAP}m \
  -Xmx${KNOTX_MAX_HEAP}m \
  -XX:ReservedCodeCacheSize=${KNOTX_CODE_CACHE}m \
  -XX:+UseBiasedLocking \
  -XX:BiasedLockingStartupDelay=0 \
  ${KNOTX_JMX_OPTS} \
  ${KNOTX_DEBUG_OPTS} \
  ${KNOTX_GC_OPTS} \
  ${KNOTX_EXTRA_OPTS} \
  -Dio.knotx.KnotxServer.options.config.serverOptions.port=${KNOTX_PORT} \
  -cp "<%= @knotx_root_dir %>/app/*" io.knotx.launcher.LogbackLauncher \
  ${KNOTX_CONFIG} \
  ${KNOTX_CONFIG_EXTRA} >> <%= @knotx_log_dir %>/<%= @knotx_id %>.log 2>&1'
Type=simple
PIDFile=/var/run/<%= @knotx_id %>
User=<%= @knotx_user %>

[Install]
WantedBy=default.target
